# com.liaoyijie.eventscope.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liaoyijie.eventscope</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/miniconda3/bin/streamlit</string>
        <string>run</string>
        <string>/Users/liaoyijie/金融科技/app.py</string>
        <string>--server.headless</string>
        <string>true</string>
        <string>--server.port</string>
        <string>8505</string>
        <string>--server.address</string>
        <string>0.0.0.0</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/liaoyijie/金融科技</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/liaoyijie/金融科技/streamlit.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/liaoyijie/金融科技/streamlit.log</string>
</dict>
</plist>
